# ansible/deploy/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/home/$user$/deploy/current/env3/run/gunicorn.sock

[Install]
WantedBy=sockets.target