# infra/playbooks/dev/yarn_devw.service
[Unit]
Description=Yarn devw service
After=yarn.service

[Service]
User=vagrant
Group=vagrant
TimeoutStartSec=900
WorkingDirectory=/vagrant/zapisy/
ExecStart=/usr/bin/yarnpkg dev:watch
Type=simple

[Install]
WantedBy=vagrant.mount vagrant-zapisy-node_modules.mount