# infra/playbooks/templates/django_cleanup.timer
[Unit]
Description=Clear old django database sessions

[Timer]
OnCalendar=FRI *-*-* 23:00:00
Persistent=true

[Install]
WantedBy=timers.target
